# docs/systemd/vault-snapshot-agent.service
# /etc/systemd/system/vault-snapshot-agent.service
[Unit]
Description=Vault Snapshot Agent

[Service]
ExecStart=/path/to/vault-snapshot-agent/
